# docs/extraction-server.service
# systemd-Dienst fuer den dedizierten Server.
# Auf dem VPS ablegen als:  /etc/systemd/system/extraction-server.service
# Danach:  sudo systemctl daemon-reload && sudo systemctl enable --now extraction-server
[Unit]
Description=Extraction Shooter - dedizierter Server
After=network.target

[Service]
Type=simple
User=extraction
WorkingDirectory=/opt/extraction
ExecStart=/opt/extraction/godot/Godot_v4.7.1-stable_linux.x86_64 --headless --path /opt/extraction/projekt -- --server --port 24567
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
